# bosting-cp/templates/freebsd/chef-client.erb
#!/bin/sh

# PROVIDE: chef-client
# REQUIRE: LOGIN cleanvar
# KEYWORD: shutdown

. /etc/rc.subr

name=chef_client
rcvar=chef_client_enable

pidfile="/var/run/${name}.pid"
command="<%= @command %>"
command_args="--daemonize --runlist recipe[bosting-generator] --pid ${pidfile}"
procname="<%= @procname %>"

load_rc_config $name
run_rc_command "$1"
